# debian/bionet2-data-manager.postrm
#!/bin/bash

BDM_DB_FILE=/var/lib/bionet-data-manager/bdm.db

if [ "$1" = "purge" ]; then
    update-rc.d bionet-data-manager remove

    rm -f $BDM_DB_FILE
fi


