# debian/bionet2-data-manager.postrm
#!/bin/bash

if [ "$1" = "purge" ]; then
    update-rc.d bionet-data-manager remove
fi


