# debian/bionet2-data-manager.prerm
#!/bin/bash

BDM_DB_FILE=/var/lib/bionet-data-manager/bdm.db

if [ "$1" = "remove" ]; then
    invoke-rc.d bionet-data-manager stop
    rm -f $BDM_DB_FILE
fi


